# etc/rc.local
#!/bin/bash -e
#
# rc.local for raspimjpeg
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#

sudo mkdir -p /dev/shm/mjpeg
sudo su -c 'raspimjpeg > /dev/null &' 

cd /home/pi/mjpg-streamer
mjpg_streamer -i "input_file.so -d 0.05 -f /dev/shm/mjpeg -n cam.jpg" -o "output_http.so -w ./www -p 8080"&

exit 0
